# vps/rowboat/crontab
# New Coworker VPS Scheduled Tasks
# Install: crontab /opt/rowboat/crontab

# Heartbeat check every 2 minutes
*/2 * * * * /opt/newcoworker/scripts/heartbeat.sh >> /var/log/heartbeat.log 2>&1

# Rotate rowboat logs daily at midnight
0 0 * * * docker compose -f /opt/rowboat/docker-compose.yml exec rowboat /usr/local/bin/rotate-logs.sh >> /var/log/rowboat-rotate.log 2>&1

# Reboot weekly on Sunday 3am (maintenance window)
0 3 * * 0 /sbin/reboot
